# installations/etc/supervisor/conf.d/001_APPNAME.conf
[program:$APPNAME]
command=/home/sites/virtualenvs/$APPNAME/bin/gunicorn_start 
user=sites
stdout_logfile=/home/sites/virtualenvs/$APPNAME/logs/gunicorn_supervisor_stdout.log
stderr_logfile=/home/sites/virtualenvs/$APPNAME/logs/gunicorn_supervisor_stderr.log
redirect_stdrr=true
autostart=true
autorestart=true
