# installations/etc/supervisor/conf.d/003_APPNAME.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat_sportspal]
; Set full path to celery program if using virtualenv
command=/home/sites/virtualenvs/$APPNAME/bin/celerybeat_start

; remove the -A myapp argument if you are not using an app instance

directory=/home/sites/apps/$APPNAME
user=sites
numprocs=1
stdout_logfile=/home/sites/virtualenvs/$APPNAME/logs/celery_beat_stdout.log
stderr_logfile=/home/sites/virtualenvs/$APPNAME/logs/celery_beat_sterr.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
